# librarydown-api.service
[Unit]
Description=LibraryDown Social Media Downloader API
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/librarydown
Environment="PATH=/root/librarydown/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/root/librarydown/.env

# Start FastAPI server on port 8001 (different from apichecker)
ExecStart=/root/librarydown/venv/bin/uvicorn src.api.main:app --host 0.0.0.0 --port 8001 --workers 1

# Resource limits for Termux environment (961MB total)
MemoryLimit=200M
CPUQuota=40%
TasksMax=50

# Process limits
LimitNOFILE=1024
LimitNPROC=32

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=librarydown-api

[Install]
WantedBy=multi-user.target